# python-producer/supervisord.conf
[supervisord]
nodaemon=true

[program:generator]
command=python /kafka-python-0.9.5/generator.py
